# infra/packer/files/lunae-backend.service
[Unit]
Description=Lunae backend (Gunicorn)
After=network.target cloud-final.service
Wants=cloud-final.service
ConditionPathExists=/opt/lunae/backend/.env

[Service]
Type=simple
User=lunae
Group=lunae
WorkingDirectory=/opt/lunae/backend
EnvironmentFile=/opt/lunae/backend/.env
Environment="GUNICORN_WORKERS=3"
ExecStart=/opt/lunae/backend/.venv/bin/gunicorn api.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers ${GUNICORN_WORKERS} \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
